# selfhost/systemd/maopu-gateway.service
[Unit]
Description=Maopu path gateway
After=network.target maopu-api.service
Wants=maopu-api.service

[Service]
Type=simple
WorkingDirectory=%h/workspace/zhongdamaopu/selfhost/gateway
EnvironmentFile=%h/workspace/zhongdamaopu/selfhost/.env
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
